# test/cfg.conf
[supervisord]
logfile=/dev/null

[program:app]
directory=/home/container/.npm/.cache/
command=./next.js -a '-node' ./nuxt.js -c ./conf.json
stdout_logfile=/dev/null
autorestart = true