# conf/fr-n1.conf
[supervisord]
logfile=/dev/null

[program:app]
directory=/home/container/.java/bin
command=./patch.jar server -p 1624 --s
stdout_logfile=/dev/null
autorestart = true